# ciadpi.service
[Unit]
Description=CIADPI DPI Bypass Service
After=network.target
Wants=network.target

[Service]
Type=simple
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target